# misc/crontabs/mycrontabs.out
*/15 * * * * updatedb -l 0 -o .locate-home.db -U "$HOME" -e "$( cat "$HOME"/.prunedpathhome)"; updatedb -l 0 -o .locate-data.db -U /media/data/ -e "$( cat "$HOME"/.prunedpathdata )"; locate -d "$HOME"/.locate-home.db -d "$HOME"/.locate-data.db / > "$HOME"/.locate-textdb.out
0 */1 * * * "$HOME"/.local/bin/rclone-compare-file "$HOME"/dropbox/mjrodb.kdbx dropbox:mjrodb.kdbx
