# dot-keepassxc/.config/systemd/user/keepassxc.service
[Unit]
Description=KeepassXC background daemon

[Service]
ExecStart=/usr/bin/keepassxc

[Install]
WantedBy=default.target
