# systemd/etc/systemd/system/suspend@.service
[Unit]
Description=User suspend actions
Before=sleep.target

[Service]
User=%I
Type=forking
Environment=DISPLAY=:0
ExecStartPre=
ExecStart=/usr/bin/xscreensaver-command -lock
ExecStartPost=/usr/bin/sleep 1

[Install]
WantedBy=sleep.target
